# layers/macos_arm64/etc/init.d/ovmservice
#!/sbin/openrc-run
name="ovm oneshut service scripts"

depend() {
	before podman sshd
}

start_pre() {
	modprobe virtiofs || {
		echo "Error: virtiofs kernel module load failed"
		exit 100
	}
}

format_and_mount_var() {
	if [[ ! -b /dev/vdb ]]; then
		echo "/dev/vdb not find"
		exit 100
	fi

	mount_param="data=ordered"
	disk_uuid=8c29bbdf-db2c-4ef6-ad7f-76bd39481d0f

	blkid -U $disk_uuid
	ret=$?
	if [[ $ret -ne 0 ]]; then
		# initialize external disk
		mkfs.ext4 -U $disk_uuid -F /dev/vdb || {
			echo "Format /dev/vdb with $disk_uuid failed"
			exit 100
		}
		mkdir -p /dev/shm/var_temp
		mount -o "$mount_param" /dev/vdb /dev/shm/var_temp || {
			echo "Mount /dev/vdb into /dev/shm/var_temp failed"
			exit 100
		}
		cp -rf /var/* /dev/shm/var_temp || {
			echo 'Copy /var/ into /dev/shm/var_temp/ failed'
			exit 100
		}
		sync
		umount /dev/shm/var_temp || {
			echo "umount /dev/shm/var_temp failed"
			exit 100
		}
		mount -o "$mount_param" /dev/vdb /var/ || {
			echo "Mount $disk_uuid into /var failed"
			exit 100
		}
	else
		mount -o "$mount_param" /dev/vdb /var/ || {
			echo "Mount $disk_uuid into /var failed"
			exit 100
		}
	fi
}

mount_initfs() {
	ign_sh="/tmp/initfs/ovm_ign.sh"
	initfs_uuid=c5c8073159f40aa69d83a1e6c7aafb16b1e5
	initfs_mount_point="/tmp/initfs"
	mount -m -t virtiofs "${initfs_uuid}" "$initfs_mount_point" || {
		echo " Error: Mount virtiofs /tmp/initfs failed"
		exit 100
	}

	if [[ -f "${ign_sh}" ]]; then
		bash +x "${ign_sh}" || {
			echo "Error: Ignition scripts exec failed"
			exit 100
		}
	else
		echo "Error: ${ign_sh} not find"
		exit 100
	fi
}

umount_efi() {
	# if we do not unmount /boot/efi, when the vm got killed, the fsck will
	# take 1 second to fix the /boot/efi filesystem.
	# This is a complete waste of time and unnecessary
	umount /boot/efi/
}

# Update grub only running at first boot
update_grub() {
	if [[ -f /.first_boot ]]; then
		update-grub && rm /.first_boot
		sync
	fi
}

# Delete any proxy environment
delete_proxy_config() {
	rm -rf /etc/profile.d/proxy.sh
}

start() {
	delete_proxy_config
	# update_grub happend before umount_efi
	update_grub
	umount_efi
	# Format /dev/vdb with $uuid first
	format_and_mount_var
	# Mount /tmp/initfs(virtioFS), and exec $ign_sh
	mount_initfs
	sync
}
