# layers/amd64_wsl2/tree.list
./
./opt
./opt/logger_wrapper.sh
./opt/ovmd
./etc
./etc/init.d
./etc/init.d/podman
./etc/init.d/kernelogger
./etc/conf.d
./etc/conf.d/podman
./etc/runlevels
./etc/runlevels/stop_all
./tree.list
